# SW-Dev/Encoder/encoder.service
[Unit]
Description= Encoder Service

[Service]
Type=oneshot
ExecStart=/usr/local/scripts/counter.sh

[Install]
WantedBy=multi-user.target
